# configs/weird.conf
[program:weird]
command=~/facct-tutorial/gunicorn_start
user=username
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=~/facct-tutorial/logs/gunicorn-error.log